# server/pkg/com.m4server.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.m4server.daemon</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/m4server</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>5</integer>

  <key>StandardOutPath</key>
  <string>/Library/Logs/M4Server/stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Library/Logs/M4Server/stderr.log</string>
</dict>
</plist>
